# deploy/systemd/openclaw-gateway-restart.service
[Unit]
Description=Restart OpenClaw Gateway

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user restart openclaw-gateway.service
